# supervisord/supervisord.conf
[supervisord]
nodaemon=true
[program:reverse_proxy]
command=/usr/sbin/nginx -g "daemon off;"
[program:rest_flask]
directory=/app
command=/usr/local/bin/python3 app.py

